# apps/drama/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:crawler]
command=/opt/app/crawler.py

[program:grpc_servicer]
command=/opt/app/grpc_servicer.py
